# roles/oai-epc/files/mme.service
[Unit]
Description=open air interface epc mme

[Service]
ExecStart=/root/openair-cn/SCRIPTS/run_mme

[Install]
WantedBy=multi-user.target