# roles/oai-epc/files/spgw.service
[Unit]
Description=open air interface epc spgw

[Service]
ExecStart=/root/openair-cn/SCRIPTS/run_spgw

[Install]
WantedBy=multi-user.target